# dotfiles/autostart/.config/autostart/org.gnome.terminal.desktop.desktop
[Desktop Entry]
Name=Terminal
Exec=gnome-terminal
Type=Application
Name[en_GB]=org.gnome.terminal.desktop
